# coral/services/command_handler.service
[Unit]
Description=Sapling serial command handler
StartLimitIntervalSec=0 #TODO make more robust
[Service]
Type=simple
Restart=always
RestartSec=1
User=mendel
ExecStart=/bin/sh /home/mendel/src/run.sh

[Install]
WantedBy=multi-user.target